# Scripts/neobot.conf
description "start and stop the neobot server"
version "1.0"
author "neobot"

start on runlevel [2345]
stop on runlevel [06]

script
exec /home/ubuntu/git/bin/NServer
end script
